# resources/supervisord.conf
[supervisord]
nodaemon=true

[program:ipython]
command=ipython notebook --notebook-dir=/notebooks
